# bash/profile.d/50-bashrc.sh
if [ "$BASH_VERSION" -a -f ~/.bashrc ]; then
    . ~/.bashrc
fi
