# bash/profile.d/bashrc.sh
if [ "$BASH" ]; then
    # Bash script
    if [[ ! -o posix && -f ~/.bashrc ]]; then
        . ~/.bashrc
    fi
fi
